# supervisord/ollama.conf
[program:ollama]
command=/usr/local/bin/ollama serve
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=1
user=ubuntu
environment=
    HOME="/home/ubuntu",
    OLLAMA_HOST="0.0.0.0:11434",
    OLLAMA_KEEP_ALIVE="10m",
    OLLAMA_MAX_LOADED_MODELS="1",
    OLLAMA_NUM_PARALLEL="1",
    NVIDIA_VISIBLE_DEVICES="all",
    NVIDIA_DRIVER_CAPABILITIES="all"
stdout_logfile=/var/log/ollama.log
stderr_logfile=/var/log/ollama.err
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile_backups=3
