# supervisord/open-webui.conf
[program:open-webui]
command=/usr/local/bin/open-webui serve --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=2
user=ubuntu
environment=
    HOME="/home/ubuntu",
    DATA_DIR="/home/ubuntu/.local/share/open-webui",
    OLLAMA_BASE_URL="http://127.0.0.1:11434",
    WEBUI_AUTH="true",
    ENABLE_MEMORY="true",
    MEMORY_RETRIEVAL_TOP_K="5",
    RAG_EMBEDDING_ENGINE="ollama",
    RAG_EMBEDDING_MODEL="nomic-embed-text",
    ENABLE_RAG_WEB_SEARCH="true",
    RAG_WEB_SEARCH_ENGINE="duckduckgo",
    DO_NOT_TRACK="true",
    ANONYMIZED_TELEMETRY="false",
    WEBUI_SECRET_KEY="%(ENV_WEBUI_SECRET_KEY)s"
stdout_logfile=/var/log/open-webui.log
stderr_logfile=/var/log/open-webui.err
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
